# autostart/qubes-app-menu.desktop
[Desktop Entry]
Type=Application
Exec=/usr/libexec/qubes-app-menu-start
Icon=qubes-manager
Terminal=false
Name=Qubes Application Menu
GenericName=Qubes Application Menu
StartupNotify=false
Hidden=false
NotShowIn=X-QUBES
